# deploy/d-brain-process.service
[Unit]
Description=d-brain Daily Processing

[Service]
Type=oneshot
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/projects/agent-second-brain
ExecStart=/home/YOUR_USER/projects/agent-second-brain/scripts/process.sh
Environment=PYTHONUNBUFFERED=1
